# docker/nginx/supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true

[program:crond]
command=/usr/sbin/crond -f -c /var/spool/cron/crontabs -L /proc/1/fd/2 -l 0
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
